# com.showcrewnetwork.login.service
[Unit]
Description=Automated login script for showcrewnetwork.com

[Service]
WorkingDirectory=~/login-main
ExecStart=/bin/bash -c 'source ~/login-main/venv/bin/activate && ~/login-main/venv/bin/python ~/login-main/login_script.py'
Environment=USERNAME=
Environment=PASSWORD=

[Install]
WantedBy=multi-user.target
